# infrastructure/preview/node-monitor/flannel-watchdog.timer
[Unit]
Description=Run flannel-watchdog every 30 seconds
Documentation=https://github.com/danify-cz/danify-app/issues/127

[Timer]
# First run 30 s after boot, then every 30 s thereafter.
OnBootSec=30s
OnUnitActiveSec=30s
AccuracySec=5s
Unit=flannel-watchdog.service

[Install]
WantedBy=timers.target
